# etc/init/reverse-tether.conf
start on filesystem and started docker
stop on runlevel [!2345]
respawn

chdir /opt/reverse-tether
exec su - root -c /bin/bash -c "/opt/reverse-tether/upstart.sh"
